# build/docker/crontab
#每分钟运行trigger
* * * * * /usr/local/php/bin/php /var/www/html/index.php schedule_task task_trigger>>/var/www/logs/task.log

